# image/mbr/hybrid-raid-layout/templates/rootfs-overlay/etc/systemd/system/disk-expansion.service.d/override.conf.in
[Service]
# Override environment for hybrid RAID layout using stable device IDs
Environment=SSD_IDS=__SSD_IDS__
