# layer/storage/luks/device/rootfs-overlay/etc/profile.d/luks-status.sh
# Show LUKS status on login
if [ -x /usr/local/bin/luks-status ]; then
  echo "🔒 LUKS Encryption Status:"
  /usr/local/bin/luks-status | head -10 | tail -8
  echo
fi
